# root/etc/periodic/daily/certbot
#!/bin/sh

echo "Running daily Certbot renewal check"
exec /usr/bin/certbot renew --verbose --nginx \
	--non-interactive \
	--user-agent-comment 'bookstack-certbot'

